# deploy/github_ip_update.service
[Unit]
Description=Update GitHub Webhook IPs

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update_github_ips.sh
